# examples/com.example.wc-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Exemple launchd macOS. Remplacez les chemins et les variables, puis :
     cp examples/com.example.wc-tracker.plist ~/Library/LaunchAgents/
     launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.example.wc-tracker.plist
     Décharger : launchctl bootout gui/$(id -u)/com.example.wc-tracker -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.wc-tracker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/ABSOLUTE/PATH/TO/senegal_wc_tracker.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>WC_TRACK_TEAM</key>
        <string>Senegal</string>
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>PUT_YOUR_BOT_TOKEN</string>
        <key>TELEGRAM_CHAT_ID</key>
        <string>PUT_YOUR_CHAT_ID</string>
        <key>FOOTBALL_API_KEY</key>
        <string>PUT_YOUR_API_SPORTS_KEY</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
